# meta-platinenmacher/recipes-core/persist-config/files/persist-config.service
[Unit]
Description=Persist system configuration to userdata partition
After=multi-user.target

[Service]
Type=oneshot
ExecStart=/usr/sbin/persist-config.sh
StandardOutput=journal

[Install]
WantedBy=multi-user.target
