# meta-platinenmacher/recipes-core/persist-config/files/persist-config.service
[Unit]
Description=Persist system configuration to userdata partition
After=local-fs.target restore-config.service
Wants=restore-config.service

[Service]
Type=simple
ExecStart=/usr/sbin/persist-config.sh daemon
ExecStop=/usr/sbin/persist-config.sh persist
Restart=on-failure
RestartSec=10
StandardOutput=journal
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
