# meta-platinenmacher/recipes-core/thebrutzler-user/files/serial-getty@.service
[Service]
ExecStart=
ExecStart=-/sbin/agetty -o '-p -f thebrutzler' --noclear --autologin thebrutzler %I $TERM
